# docker/etc/s6/fastgithub/run
#!/bin/bash

# 启动FastGithub加速
exec /opt/fastgithub/fastgithub
